# supervisord.conf
[supervisord]
nodaemon=true

[inet_http_server]
port=0.0.0.0:1001
username=admin
password=q1w2e3@123

[program:elasticpotpy]
directory=/opt/ElasticpotPY/
command=/bin/sh -c "python3 main.py"